# etc/rc.local
#! /bin/bash

# auto video stream
/root/mjpgstream/beginstream.sh

# Attempt to start the Anubis background Service
/usr/local/bin/anubiscon start

exit 0
