# etc/rc.local
#! /bin/bash

# Attempt to start the Anubis video streamer
/usr/local/bin/videocon start &

# Attempt to start the Anubis background Service
/usr/local/bin/anubiscon start &

exit 0
